# provision/ansible/roles/podman/templates/etc/systemd/podman.socket.d/override.conf.j2
[Socket]
SocketMode={{ podman_socket_mode }}
SocketUser={{ podman_socket_user }}
SocketGroup={{ podman_socket_group_name }}
